# root/etc/s6-overlay/s6-rc.d/svc-trojan-go/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec /usr/sbin/trojan-go -config /etc/trojan-go/config.json
